# systemd/user/fehbg.service
[Unit]
Description=Background Setter

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/home/jon/.config/systemd/user/scripts/fehbg
