# systemd/user/fehbg.service
[Unit]
Description=Background Setter

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/home/jon/bin/fehbg
