# systemd/user/polipo.service
[Unit]
Description=Polipo Proxy Server

[Service]
ExecStart=/usr/bin/polipo -c ~/.config/polipo/config

[Install]
WantedBy=default.target
