# systemd/user/redshift@.service
[Unit]
Description=Redshift Per-Screen

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/redshift -v -m randr:crtc=%i
Restart=on-failure
