# systemd/user/unclutter.service
[Unit]
Description=Mouse Hider

[Service]
ExecStart=/usr/bin/unclutter -display :0 -idle 3
Restart=on-failure
