# systemd/user/updates.timer
[Unit]
Description=Inform about available updates hourly

[Timer]
OnCalendar=*:19
Persistent=true

[Install]
WantedBy=default.target
