# systemd/user/updates.service
[Unit]
Description=Inform about available updates

[Service]
ExecStart=/home/jon/.config/systemd/user/scripts/updates
