# systemd/user/urxvtd.service
[Unit]
Description=urxvt Daemon

[Service]
Type=simple
ExecStart=/usr/bin/urxvtd
