# systemd/user/urxvtd.service
[Unit]
Description=urxvt Daemon

[Service]
ExecStart=/usr/bin/urxvtd
Restart=on-failure
