# systemd/user/xresources.service
[Unit]
Description=Xresources Settings
ConditionPathExists=%h/.Xresources

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/xrdb -merge %h/.Xresources
